# g5-exporters/systemd-user/docker-stats-metrics.service
[Unit]
Description=Emit Docker stats Prometheus textfile metrics

[Service]
Type=oneshot
ExecStart=%h/homelab-exporters/scripts/docker-stats-metrics.sh
Environment=TEXTFILE_DIR=%h/homelab-exporters/textfile
